# TARGET_PACKAGES/binary/startupscript/odroid.service
[Unit]
Description=set ip and mount partition

[Service]
Type=oneshot
ExecStart=/etc/run_odroid_service

[Install]
WantedBy=default.target
